# etc/udev/rules.d/99-monitor-hotplug.rules
KERNEL=="card0", ACTION=="change", RUN+="/opt/udev-monitor-hotplug/usr/local/bin/monitor-hotplug.sh"
